# sgai-auto-shutdown.service
[Unit]
Description=SGAI Auto-Shutdown Monitor
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/workspace/ragflow
ExecStart=/usr/bin/python3 /home/ubuntu/workspace/ragflow/auto_shutdown_monitor.py
Restart=always
RestartSec=10
# Nessun limite: dopo un crash o reboot Docker, il monitor riparte sempre
StartLimitIntervalSec=0
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

# Verifica deploy: systemctl is-active sgai-auto-shutdown.service (atteso: active)
